# deploy/macos/com.yanxi.route-warden.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yanxi.route-warden</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yanxi/.cargo/bin/route-warden</string>
    <string>--config</string>
    <string>/Users/yanxi/.route-warden/config.yaml</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/yanxi</string>

  <key>StandardOutPath</key>
  <string>/tmp/route-warden.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/route-warden.err.log</string>
</dict>
</plist>
